# ntpdate-temp/0ntpdate.cron
#!/bin/bash
exec /sbin/ntpdate ntp1.atto.io >&2> /dev/null
